# deploy/signal-api.service
[Unit]
Description=KOSPI Signal API (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/apps/kospi-scanner/signal-api
Environment="SIGNAL_API_DATA_DIR=/opt/apps/kospi-scanner/data"
Environment="SIGNAL_API_CORS_ORIGINS=https://sigbora.com"
ExecStart=/opt/apps/kospi-scanner/.venv/bin/uvicorn app.main:app \
          --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
